# deploy/systemd/signal-inbound-collector.service
[Unit]
Description=SigPro Signal Inbound Collector (JSONL -> SQLite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# User inherited from systemd --user manager
WorkingDirectory=/home/james/.openclaw/workspace-sigpro
ExecStart=/usr/bin/env python3 /home/james/.openclaw/workspace-sigpro/scripts/signal_inbound_collector.py --follow --poll-ms 1500
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/james/.openclaw/workspace-sigpro/.openclaw /home/james/.local/share/signal-cli/attachments

[Install]
WantedBy=default.target
